# systemd/system/zypper_download.service
[Unit]
Description=Download Zypper Updates in the background
Type=exec

[Service]
ExecStart=zypper-dup-download
