# systemd/user/jellyfin.service
[Unit]
Description=Jellyfin

[Service]
ExecStart=/home/henry/.local/software/jellyfin/jellyfin

[Install]
WantedBy=default.target
