# systemd/user/local-web.service
[Unit]
Description=Local Web

[Service]
ExecStart=/home/henry/.local/bin/local-web -a $OPEN_WEATHER_API_KEY

[Install]
WantedBy=default.target
