# systemd/user/wallpaper.timer
[Unit]
Description=Change Wallpaper Periodically
Requires=wallpaper.service

[Timer]
OnCalendar=00/1:00
Persistent=true
AccuracySec=1s
Unit=wallpaper.service

[Install]
WantedBy=timers.target
