# systemd/user/wallpaper.service
[Unit]
Description=Change Wallpaper Periodically

[Service]
Type=oneshot
ExecStart=timed_wallpaper.py $HOME/Pictures/background/

